# debian/cron.d
*/5 * * * *     munin if [ -x /usr/bin/munin-cron ]; then /usr/bin/munin-cron --config /etc/munin-json/munin.conf; fi

